# scripts/stillestille.service
[Unit]
Description=HDD Standby

[Service]
ExecStart=/usr/bin/stillestille.sh

[Install]
WantedBy=multi-user.target 
